# systemd/easyape.service
[Unit]
Description=EasyApe - Text to Stake
After=network.target

[Service]
Type=simple
User=root

# ✅ MUST match your repo location
WorkingDirectory=/root/TEASY-test

ExecStart=/root/TEASY-test/.venv/bin/python \
    -u -m stakechat_bot.main \
    --config /root/TEASY-test/config.yaml

Restart=always
RestartSec=5

# Clean shutdown
TimeoutStopSec=10
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
